# docker/production/supervisor.d/laravel-queue.ini
[program:laravel-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/backend/artisan queue:work --tries=3
autostart=true
autorestart=true
numprocs=5
redirect_stderr=true
stdout_logfile=/var/www/backend/storage/logs/worker.log
